# opt/ncubed/updates/20220419-1344/ncubed/oneshot.service/ncubed-oneshot.service
[Unit]
Description=[ NCUBED ] oneshot scripts

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/ncubed/oneshot.service/oneshot

[Install]
WantedBy=multi-user.target
